# deploy/homebot-kiosk.service
[Unit]
Description=HomeBot Chromium Touch Dashboard (fullscreen kiosk)
After=homebot-server.service graphical-session.target
Requires=homebot-server.service

[Service]
Type=simple
EnvironmentFile=%h/homebot/deploy/env
Environment=DISPLAY=:0
ExecStart=%h/homebot/deploy/launch-kiosk.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
